# systemd/zdns-rest.service
[Unit]
Description=ZDNS REST API Server
Documentation=https://github.com/smeinecke/zdns-rest
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/zdns-rest server --config /etc/zdns-rest/zdns-rest.conf
Restart=on-failure
RestartSec=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
